# system/systemd/hidloom-network-late.timer
[Unit]
Description=Start networking outside the keyboard boot path

[Timer]
OnBootSec=25s
AccuracySec=1s
Unit=hidloom-network-late.service

[Install]
WantedBy=timers.target
